# deploy/systemd/netfirewall-daemon.service
[Unit]
Description=NetFirewall privileged daemon (network mutations)
Documentation=file:///opt/netfirewall/README.txt
After=network-online.target postgresql.service
Wants=network-online.target

# The daemon owns OS mutations. The Web cannot start cleanly without it.
Before=netfirewall-web.service

[Service]
Type=notify
NotifyAccess=main

# Runs as root because it has to write /etc/network/*, /etc/NetworkManager/*,
# /etc/netplan/* and invoke nmcli / ip / ifup. CapabilityBoundingSet below
# clamps that root down to the strict minimum.
User=root
Group=root

WorkingDirectory=/opt/netfirewall/daemon
ExecStart=/opt/netfirewall/daemon/netfirewall-daemon
EnvironmentFile=-/etc/netfirewall/daemon.env

# /run/netfirewall is created by tmpfiles.d (deploy/tmpfiles/netfirewall.conf
# → /etc/tmpfiles.d/netfirewall.conf), NOT by RuntimeDirectory=. Reason:
# systemd doesn't expose RuntimeDirectoryGroup=, so RuntimeDirectory= would
# leave the dir as root:root and break the Web's ability to traverse to the
# socket. tmpfiles.d lets us declare root:netfirewall 0750 directly.

# Persistent state (logs). Owned by root since the daemon runs as root.
StateDirectory=netfirewall/daemon
StateDirectoryMode=0750
LogsDirectory=netfirewall/daemon
LogsDirectoryMode=0750

# ─── Hardening ─────────────────────────────────────────────────────────
# Capability clamp — only what's needed:
#   CAP_NET_ADMIN   — modify interfaces, routes, addresses
#   CAP_DAC_OVERRIDE — read/write files in /etc/network /etc/NetworkManager
#   CAP_NET_RAW     — required by some `ip` / `nft` operations
#   CAP_CHOWN       — chown the control socket to root:netfirewall so the
#                     Web user (in that group) can connect over it
# Everything else (including CAP_SYS_ADMIN, CAP_SYS_MODULE, CAP_SYS_PTRACE,
# CAP_AUDIT_WRITE, CAP_KILL, etc.) is unavailable.
CapabilityBoundingSet=CAP_NET_ADMIN CAP_DAC_OVERRIDE CAP_NET_RAW CAP_CHOWN
AmbientCapabilities=

NoNewPrivileges=yes

# ProtectSystem=strict makes /usr /boot /efi read-only for the process.
# ReadWritePaths re-opens only the directories the daemon legitimately writes.
# The `-` prefix on a path makes it optional — Debian without NetworkManager,
# Rocky without netplan, etc. don't have all these dirs, and a missing one
# fails the mount namespace setup with status=226/NAMESPACE. The daemon
# detects which network stack is present at runtime; we just give the
# namespace permission to write to any of them that happen to exist.
ProtectSystem=strict
ReadWritePaths=-/etc/network -/etc/NetworkManager -/etc/netplan -/etc/nftables.conf /var/lib/netfirewall /var/log/netfirewall
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
SystemCallArchitectures=native
# @system-service is a curated allow-list; we additionally deny @mount, @swap,
# @reboot, @raw-io which the daemon never needs.
SystemCallFilter=@system-service
SystemCallFilter=~@mount @swap @reboot @raw-io @cpu-emulation

# NOTE: MemoryDenyWriteExecute is intentionally omitted — the .NET JIT writes
# executable pages and the runtime would crash. Switch to ReadyToRun / NativeAOT
# publish to enable it.

# Resource limits (sane defaults — bump if you have a busy network).
LimitNOFILE=8192
TasksMax=128

Restart=on-failure
RestartSec=5

# Stop politely; SIGINT lets ASP.NET shut down cleanly.
KillSignal=SIGINT
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
